# src/docker/root/etc/s6-overlay/s6-rc.d/frontend/run
#!/command/with-contenv bash

cd /app/frontend/dist
export NEXT_PUBLIC_API_URL=http://localhost:5000
exec s6-setuidgid abc serve -s . -l 3000
