# src/docker/root/etc/s6-overlay/s6-rc.d/frontend/run
#!/command/with-contenv bash

cd /app/frontend
exec s6-setuidgid abc npx serve -s . -l 3000
