# src/docker/root/etc/s6-overlay/s6-rc.d/backend/run
#!/command/with-contenv bash

cd /app/backend
exec s6-setuidgid abc dotnet Trailarr.Core.dll
