# src/docker/root/etc/s6-overlay/s6-rc.d/svc-backend/run
#!/command/with-contenv bash
exec s6-notifyoncheck -d -n 300 -w 1000 \
    s6-setuidgid abc \
    cd /app/build && dotnet Trailarr.Host.dll
